# audio-keepalive/audio-keepalive.service
[Unit]
Description=Audio Keepalive - Prevent audio device idle suspend
After=pipewire.service pulseaudio.service

[Service]
ExecStart=%h/.local/bin/audio-keepalive.sh
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
